# setup/ab_booting_variant/probe-update.timer
[Unit]
Description=MobileAtlas probe update timer

[Timer]
OnBootSec=1h

[Install]
WantedBy=default.target
